# Data Visualization/Bokeh/Deploy Bokeh Visualization Apps in Live Servers/Flask-with-Bokeh-Server/flask.conf
[program:flask]
command=/opt/envs/virtual/bin/gunicorn -b :8118 app:app
directory=/opt/webapps/bokehflask
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
